# sonarr-import.path
[Unit]
Description=Watch downloads folder for new TV show episodes
Documentation=https://github.com/thomasevano/sonarr-manual-import

[Path]
# Path is read from settings.json at install time
# This will be replaced by the install script
PathChanged=/path/to/downloads
Unit=sonarr-import.service

[Install]
WantedBy=multi-user.target
