# sonarr-import.path
[Unit]
Description=Watch downloads folder for new TV show episodes
Documentation=https://github.com/thomasevano/sonarr-manual-import

[Path]
# Path is read from settings.json at install time
# This will be replaced by the install script
PathChanged=/path/to/downloads
# Debounce: wait for file writes to settle (10 seconds)
TriggerLimitIntervalSec=10
TriggerLimitBurst=1

[Install]
WantedBy=multi-user.target
